# priv/install/lobber.service.d/env.conf
[Service]
Environment=DISCORD_BOT_TOKEN=""
Environment=OPENROUTER_API_KEY=""
Environment=PERPLEXITY_API_KEY=""
